# tmux-console.service
[Unit]
Description=Tmux persistent console sessions (creates console-1..5 on boot)
After=default.target

[Service]
Type=forking
ExecStart=/bin/bash %h/.tmux-persistent-console/setup.sh
ExecStop=/usr/bin/tmux kill-server
RemainAfterExit=yes

[Install]
WantedBy=default.target
